# autostart/nitrogen.desktop
[Desktop Entry]
Name=nitrogen
Icon=nitrogen
Exec=nitrogen --restore
Terminal=false
Type=Application
X-GNOME-Autostart-enabled=true
